# distribution-point/etc/crontab
#
# Crontab for __ACCOUNT__
#

0 0 * * *	$HOME/__NAME__/bin/cron/cron-midnight
0 1 * * *	$HOME/__NAME__/bin/cron/cron-daily
0 * * * *	$HOME/__NAME__/bin/cron/cron-hourly
