# distribution-point/rsync/crontab
#BEGIN-distrobox
0 * * * *	/tmp/distrobox/update-rsync
#END-distrobox
